# roles/sec/files/etc/udev/rules.d/90-yubikey.rules
SUBSYSTEM=="usb", ACTION=="remove", ENV{PRODUCT}=="1050/*", RUN+="/usr/local/bin/yubikey-lock"
