# roles/base/templates/home/.config/systemd/user/duplicity-backup.service.j2
[Unit]
Description=Duplicity backup
Wants=duplicity-backup.timer

[Service]
Type=oneshot
ExecStart=/home/{{ user_name }}/.bin/duplicity-backup.sh

[Install]
WantedBy=multi-user.target
