# roles/i3/files/etc/udev/rules.d/90-setup-displays.rules
KERNEL=="card0", SUBSYSTEM=="drm", RUN+="/usr/bin/bash -c 'sleep 5; /usr/local/bin/setup-displays'"
