# roles/i3/templates/etc/udev/rules.d/90-setup-displays.rules.j2
KERNEL=="card0", SUBSYSTEM=="drm", ENV{DISPLAY}=":0", ENV{XAUTHORITY}="/home/{{ user_name }}/.Xauthority", RUN+="/usr/bin/bash -c 'sleep 5; /usr/local/bin/setup-displays.sh'"
