# roles/i3/templates/home/.config/systemd/user/duplicity-backup.timer.j2
[Unit]
Description=Duplicity backup
Requires=duplicity-backup.service

[Timer]
Unit=duplicity-backup.service
OnCalendar={{ backup_systemd_schedule }}

[Install]
WantedBy=timers.target
